# docker-monitor.service
[Unit]
Description=Мониторинг Docker-контейнеров и Telegram-бота
After=docker.service
Requires=docker.service

[Service]
WorkingDirectory=/root/deploy_testing
ExecStart=/bin/bash /root/deploy_testing/monitor.sh
Restart=always
RestartSec=30
User=root

[Install]
WantedBy=multi-user.target
